# FinOps/crontab
00 13 * * * /home/lukestridergm/PycharmProjects/FinOps/fin-ops.sh